# deploy/systemd/extract_records_publisher.service
[Unit]
Description=NetScanner DNS Record Publisher
After=network-online.target rabbitmq-server.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/netscanner/extract_records.env
WorkingDirectory=/opt/git/purple_jo
ExecStart=/opt/git/purple_jo/venv/bin/python \
    tools/extract_records.py publish \
    --postgres-dsn ${POSTGRES_DSN} \
    --rabbitmq-url ${RABBITMQ_URL} \
    --queue-name ${QUEUE_NAME} \
    --worker-count ${PUBLISHER_WORKER_COUNT} \
    --batch-size ${PUBLISHER_BATCH_SIZE} \
    ${PUBLISHER_PURGE_FLAG} \
    ${PUBLISHER_EXTRA_ARGS}
User=netscanner
Group=netscanner

[Install]
WantedBy=multi-user.target
